# profile-dir/config-sh/profile.d/001-shell.sh
#!/bin/sh
# 001-shell.sh - included by .profile

# set shell environment
export histchars=""
export ignoreeof=0
