# profile-dir/config-sh/profile.d/005-shell.bash
#!/bin/bash
# 005-shell.bash - included by .profile

# environment and startup programs
BASH_ENV="${HOME}"/.bashrc
export BASH_ENV
