# profile-dir/config-sh/profile.d/021-vimode.sh
#!/bin/bash
# 021-vimode.sh - included by .profile

# set up VI mode
export EDITOR=vim
set -o vi
