# profile-dir/config-sh/profile.d/100-dev-rust.sh
#!/bin/sh
# 100-dev-rust.sh - included by .profile

#
# software development settings
#

# Rust
CARGO_HOME=${HOME}/.local/lib/cargo
RUSTUP_HOME=${HOME}/.local/lib/rustup
export CARGO_HOME RUSTUP_HOME
